# games-server/games-server.service
[Unit]
Description=ollij.fi realtime games server
After=network.target

[Service]
Type=simple
User=ollij
WorkingDirectory=/srv/ollij.fi/games-server
ExecStart=/usr/local/bin/node server.js
Environment=PORT=8787
# The page is served from ollij.fi (GitHub Pages); the WS is this box at
# ws.ollij.fi, so requests are cross-origin — allow the page origins.
Environment=CORS_ORIGIN=https://ollij.fi,https://www.ollij.fi
# Accounts + game history SQLite DB lives here (created on first run).
Environment=DATA_DIR=/srv/ollij.fi/games-server/data
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
